# autostart/kiosk.desktop
[Desktop Entry]
Type=Application
Name=Kiosk
Exec=/home/mila/.local/bin/kiosk.sh
X-GNOME-Autostart-enabled=true

