# autostart/kiosk.desktop
[Desktop Entry]
Type=Application
Name=Kiosk
Exec=/home/mila/menu-config/scripts/kiosk.sh
X-GNOME-Autostart-enabled=true

