# autostart/show-ip.desktop
[Desktop Entry]
Type=Application
Exec=/home/mila/menu-config/show-ip.sh
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Show IP
